# .docker/supervisor.conf
[supervisord]
nodaemon=true

[program:artisan-serve]
command=php artisan serve --host=0.0.0.0 --port=8000
autostart=true
autorestart=false
directory=/app

#[program:schedule-run]
#command=ash -c 'sleep 60 && exec php artisan schedule:run'
#autorestart=true
#directory=/app
